# Library/LaunchAgents/local.camper_van.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.camper_van</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/jay/bin/camper_van_wrapper</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/jay/.camper_van/camper_van.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jay/.camper_van/camper_van_errors.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
  </dict>
</plist>
